# logcheck.cron
# /etc/cron.d/logcheck: crontab entries for the logcheck package
MAILTO=root

@reboot   logcheck nice -n10 /usr/sbin/logcheck -R
2 * * * * logcheck nice -n10 /usr/sbin/logcheck

# vim:syn=crontab
